# systemd/titan-brain.service.template
[Unit]
Description=Titan Brain (AI Engine)
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_USER
WorkingDirectory=REPLACE_WITH_WORKDIR
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 REPLACE_WITH_WORKDIR/offchain/ml/brain.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
MemoryLimit=4G

[Install]
WantedBy=multi-user.target
