# systemd/titan-super-agent.service.template
[Unit]
Description=Titan Super Agent System
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_USER
WorkingDirectory=REPLACE_WITH_WORKDIR
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 REPLACE_WITH_WORKDIR/agents/super_agent_manager.py --mode daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
MemoryLimit=1G

[Install]
WantedBy=multi-user.target
